# adsbvue.service
[Unit]
Description=ADSb-Vue — 3D ADS-B reception viewer
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Point this at your tar1090 instance (http://127.0.0.1 if it runs on this host).
Environment=ADSB_ULTRAFEEDER=http://127.0.0.1
Environment=ADSB_PORT=24556
# Uncomment to widen the coverage envelope (all history) at the cost of load time:
# Environment=ADSB_MAX_CHUNKS=0
# Adjust the path to wherever you cloned the repo.
ExecStart=/usr/bin/python3 %h/adsbvue/server.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
